# hisobot.conf
[program:hisobot]
directory=/home/ubuntu/hisobotimbot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/hisobotimbot/logfile_err.log
stdout_logfile=/home/ubuntu/hisobotimbot/logfile.err